# mopidy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Mopidy
# Runs the Nginx daemon
# ==============================================================================
declare -a options

# Wait for Mopidy to become available
bashio::net.wait_for 4478

bashio::log.info "Starting NGinx..."
exec nginx "${options[@]}"
