# mopidy/rootfs/etc/services.d/mopidy/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Mopidy
# Take down the S6 supervision tree when Mopidy fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "Mopidy crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "Mopidy stopped, restarting..."
